# linux/bin/rc.local
#!/bin/bash
(sleep 2; /home/msd/bin/load_combined_speakers.sh) & 
# (sleep 6; echo -e 'connect EC:81:93:5A:66:BB' | bluetoothctl) &
